# profile.d/.zsh_env
# 
export PATH=$PATH:$HOME/.composer/vendor/bin

# more macOS/Bash-like word jumps
export WORDCHARS=""
